# common/wallpaper/.startup/profile/profile.d/wallpaper.sh
#!/usr/bin/env sh
export WALLPAPER=~/.wallpapers/mojave.jpg
