# macos/brew/.startup/profile/profile.d/01-brew.sh
#!/usr/bin/env sh

brew=/opt/homebrew/bin/brew
[ -x "$brew" ] && eval "$("$brew" shellenv)"

