# macos/skhd/Library/LaunchAgents/skhd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <!-- Set Label -->
  <string>skhd</string>
  <key>EnableGlobbing</key>
  <true/>
  <key>ProgramArguments</key>
  <array>
    <!-- Set command to run. -->
    <!-- Add more string keys for more arguments. -->
    <!-- Start shell as login so profile is sourced and env vars are set. -->
    <string>/usr/bin/env</string>
    <string>sh</string>
    <string>--login</string>
    <string>-c</string>
    <string>skhd</string>
  </array>

  <!-- Whether to run the program on login. -->
  <key>RunAtLoad</key>
  <true/>
  <!-- Whether to restart the program if it crashes or terminates successfully. -->
  <key>KeepAlive</key>
  <false/>

  <!-- Where to place stdout and stderr logs. -->
  <!-- Can also use /usr/local/var/log directory. -->
  <key>StandardErrorPath</key>
  <string>/tmp/skhd.err.log</string>
  <key>StandardOutPath</key>
  <string>/tmp/skhd.out.log</string>
</dict>
</plist>
